# PavillionServer/alfredo-mcp.service
[Unit]
Description=Alfredo MCP Server - Ops Interface
After=network.target

[Service]
User=pvserver
WorkingDirectory=/home/pvserver/alfredo-ops
ExecStart=/home/pvserver/alfredo-ops/venv/bin/python /home/pvserver/alfredo-ops/mcp_server.py
Restart=always
RestartSec=5
Environment=MCP_AUTH_TOKEN=AlfredoOps_Seguro_2026_!
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
